# backup.sh.timer
[Unit]
Description=Example Backup Timer

[Timer]
Persistent=false
OnCalendar=hourly
Unit=backup.sh.service

[Install]
WantedBy=default.target
